# etc/supervisor/conf.d/named.conf
[program:named]
command=/usr/sbin/bind
user=root

