# etc/supervisor/conf.d/named.conf
[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:named]
command=/usr/sbin/named
program_name=named
user=root

