# etc/supervisor/conf.d/named.conf
[program:named]
command=/usr/sbin/named
user=root

